# fastapi-csv.service
[Unit]
Description=FastAPI CSV to JSON API
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/fastapi-csv
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5s
Environment="PATH=/usr/bin:/usr/local/bin"

[Install]
WantedBy=multi-user.target
